# flights.timer
[Unit]
Description=Run flights every so often

[Timer]
OnCalendar=daily
Unit=flights.service

[Install]
WantedBy=multi-user.target
