# flights.timer
[Unit]
Description=Runs myscript every hour

[Timer]
OnCalendar=daily
Unit=flights.service

[Install]
WantedBy=multi-user.target
