# flights.service
[Unit]
Description=MyScript

[Service]
Type=simple
ExecStart=/home/ubuntu/flights/runFlights.sh
